# fleet/mesos-slave-presence@.service
[Unit]
Description=Mesos Slave Announcer
BindsTo=mesos-slave@%i.service
After=mesos-slave@%i.service

[Service]
ExecStart=/bin/bash -c " \
    while [ -n \"$IP\"]; do \
        IP=$(getent hosts %H | awk '{ print $1 }'); \
        sleep 1; \
    done; \
    while true; do \
        etcdctl set --ttl 60 /skydns/local/vogt/sky/mesos-slave-%i {\\\"host\\\":\\\"$IP\\\"} && \
        etcdctl set --ttl 60 /skydns/local/vogt/sky/mesos-slave/%i {\\\"host\\\":\\\"$IP\\\"}; \
        sleep 45; \
    done"
Restart=always
RestartSec=60s

[X-Fleet]
MachineOf=mesos-slave@%i.service
