# fleet/mnt-work.mount
[Mount]
What=/home/core
Where=/mnt/work
Type=none
Options=bind

[X-Fleet]
Global=true
