# fleet/upgrade.service
[Unit]
Description=Force Upgrade
After=fleet.service

[Service]
ExecStart=/usr/bin/update_engine_client -update

[X-Fleet]
Global=true
